# etc/systemd/system/kanidmd.service.d/override.conf
[Service]
RuntimeDirectory=kanidmd
RuntimeDirectoryMode=775
StateDirectory=kanidmd
ExecStartPost=sleep 1
ExecStartPost=/usr/bin/chmod 777 /run/kanidmd/socket